# env/windows-arm64/aws/qemu.service
[Unit]
Description=QEMU

[Service]
ExecStart=/home/ubuntu/win10-arm64.sh

[Install]
WantedBy=multi-user.target
